# crontabs
00 */8 * * * /bin/bash /usr/local/bin/logrotate
