# crontabs
15 8,17 * * * /bin/bash /usr/local/bin/logrotate
